# social-network-webapi.service
[Unit]
Description=Social Network app

[Service]
WorkingDirectory=/root/SocialNetwork.WebApi.Builded/SocialNetwork/SocialNetwork.WebApi/bin/Release/netcoreapp3.0
ExecStart=/usr/bin/dotnet /root/SocialNetwork.WebApi.Builded/SocialNetwork/SocialNetwork.WebApi/bin/Release/netcoreapp3.0/SocialNetwork.WebApi.dll
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=SocialNetwork
User=root
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
